# etc/package/preremove
#!/bin/sh

if [ -f /usr/bin/systemctl ]; then
    systemctl stop redis-stack
    systemctl disable redis-stack
fi