# etc/services/redis-stack-server.service
[Unit]
Description=redis-stack-server
After=network.target
#PartOf=redis-stack.service

[Service]
Type=simple
User=nobody
Restart=always
RestartSec=1
User=nobody
ExecStart=/opt/redis-stack/bin/redis-server /opt/redis-stack/etc/redis-stack.conf
ExecStop=/bin/kill -s TERM ${MAINPID}

[Install]
WantedBy=multi-user.target
